# dig_noise.crontab
0-59/2 * * * * env MONITOR_HOST=10.209.0.7 /home/fuzhiwen/bin/dig_noise.sh || true
